# supervisor/prod.conf
[supervisord]
nodaemon=true

[program:laravel-server]
command=php /app/artisan octane:start --host=0.0.0.0 --port=80
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=0